# deploy/com.openhinge.gateway.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openhinge.gateway</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/path/to/openhinge/dist/index.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/path/to/openhinge</string>
  <key>RunAtLoad</key>
  <true/>
  <!-- Daemon policy: restart on crash or abnormal exit, but not on clean
       shutdown (SIGTERM via `openhinge stop`). ThrottleInterval=30 prevents
       crash-loop amplification that killed OpenClaw under port exhaustion. -->
  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>
  <key>ThrottleInterval</key>
  <integer>30</integer>
  <key>StandardOutPath</key>
  <string>/path/to/openhinge/data/stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/path/to/openhinge/data/stderr.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
    <key>NODE_OPTIONS</key>
    <string>--unhandled-rejections=warn</string>
  </dict>
</dict>
</plist>
